# stage3/20-openwebrx/files/openwebrx.service
[Unit]
Description=OpenWebRX WebSDR receiver
After=multi-user.target

[Service]
Type=simple
User=openwebrx
Group=openwebrx
WorkingDirectory=/opt/openwebrx
ExecStart=/usr/bin/python3 /opt/openwebrx/openwebrx.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
